# files/etc/ansible/roles/vault/files/vault_runner.conf
[supervisord]
nodaemon=true

[program:vault]
autorestart=true
command=sudo vault server -config /etc/vault.conf
