# files/etc/ansible/roles/ansible.vault/files/vault_runner.conf
[supervisord]
nodaemon=true

[program:vault]
autorestart=true
command=vault server -config /etc/vault.conf
